# modules/ruggedpod/finish/files/etc/cron.d/i2c-config
@reboot root chown ruggedpod: /dev/i2c-*
